# idca_visualizer/test_output/sample_data.json
{
  "general": {
    "company_name": "Acme Corporation",
    "report_date": "2024-01-15",
    "tester_name": "John Doe",
    "manager_name": "Jane Smith",
    "manager_title": "Security Manager"
  },
  "test_results": {
    "total_rules": 150,
    "tested_rules": 120,
    "triggered_rules": 85
  },
  "mitre_tactics": {
    "Initial Access": {
      "name": "Initial Access",
      "test_count": 10,
      "triggered_count": 8
    },
    "Execution": {
      "name": "Execution",
      "test_count": 15,
      "triggered_count": 12
    },
    "Persistence": {
      "name": "Persistence",
      "test_count": 12,
      "triggered_count": 7
    },
    "Privilege Escalation": {
      "name": "Privilege Escalation",
      "test_count": 8,
      "triggered_count": 5
    },
    "Defense Evasion": {
      "name": "Defense Evasion",
      "test_count": 20,
      "triggered_count": 15
    },
    "Credential Access": {
      "name": "Credential Access",
      "test_count": 10,
      "triggered_count": 3
    },
    "Discovery": {
      "name": "Discovery",
      "test_count": 15,
      "triggered_count": 11
    },
    "Lateral Movement": {
      "name": "Lateral Movement",
      "test_count": 8,
      "triggered_count": 6
    },
    "Collection": {
      "name": "Collection",
      "test_count": 10,
      "triggered_count": 8
    },
    "Command and Control": {
      "name": "Command and Control",
      "test_count": 12,
      "triggered_count": 10
    }
  },
  "triggered_rules": [
    {
      "name": "Suspicious PowerShell Command",
      "mitre_id": "T1059.001",
      "tactic": "Execution",
      "confidence": 95
    },
    {
      "name": "Brute Force Attack Detected",
      "mitre_id": "T1110",
      "tactic": "Credential Access",
      "confidence": 88
    },
    {
      "name": "Process Injection",
      "mitre_id": "T1055",
      "tactic": "Defense Evasion",
      "confidence": 92
    },
    {
      "name": "Persistence via Registry",
      "mitre_id": "T1547.001",
      "tactic": "Persistence",
      "confidence": 85
    },
    {
      "name": "Data Staging",
      "mitre_id": "T1074",
      "tactic": "Collection",
      "confidence": 78
    }
  ],
  "undetected_techniques": [
    {
      "mitre_id": "T1190",
      "name": "Exploit Public-Facing Application",
      "tactic": "Initial Access",
      "criticality": "Critical"
    },
    {
      "mitre_id": "T1053",
      "name": "Scheduled Task/Job",
      "tactic": "Persistence",
      "criticality": "High"
    },
    {
      "mitre_id": "T1003",
      "name": "OS Credential Dumping",
      "tactic": "Credential Access",
      "criticality": "Critical"
    }
  ],
  "recommendations": [
    {
      "category": "Technical",
      "text": "Implement PowerShell script block logging and enhanced monitoring"
    },
    {
      "category": "Process",
      "text": "Establish regular security assessment schedule for new attack patterns"
    },
    {
      "category": "Technical",
      "text": "Deploy EDR solution for better visibility into process injection attempts"
    },
    {
      "category": "Training",
      "text": "Conduct security awareness training focusing on phishing attacks"
    }
  ]
}